# docker/supervisor/app.ini
[program:nginx]
command=/usr/sbin/nginx
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true